# conf/crontab
# crontab

0 0 * * *		cd /home/zeek/MyProject/auto-check-in/run && ./auto_check_in.sh rq
